# recipes-core/systemd/systemd-conf/80-etnaviv-devcoredump.rules
ACTION=="add", SUBSYSTEM=="devcoredump", RUN+="/usr/bin/etnaviv-devcoredump /sys/%p"
